# src/web/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
command=uwsgi --ini /srv/uwsgi.ini

[program:nginx]
command=nginx
